# etc/cron.d/s3backup
PATH = '/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:'
0 3 * * *   root /usr/local/s3backup/backup.py 2>&1
